# templates/database/mariadb/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql --init-file=/etc/mysql/init.sql
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mariadb.log
stderr_logfile=/var/log/supervisor/mariadb-err.log

[program:mariadb-init]
command=/bin/bash -c "while ! mysqladmin ping --silent 2>/dev/null; do sleep 1; done; mysql -u root -proot < /usr/share/phpmyadmin/sql/create_tables.sql; mysql -u root -proot -e \"GRANT ALL PRIVILEGES ON phpmyadmin.* TO 'root'@'localhost'; FLUSH PRIVILEGES;\""
autostart=true
autorestart=false
startsecs=0
stdout_logfile=/var/log/supervisor/mariadb-init.log
stderr_logfile=/var/log/supervisor/mariadb-init-err.log

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/apache2.log
stderr_logfile=/var/log/supervisor/apache2-err.log
